# packaging/systemd/oadtd@.service
[Unit]
Description=Open Agentic Threat Defense instance %i
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=oadtd
Group=oadtd
EnvironmentFile=-/etc/oadtd/%i.env
Environment=OATD_ADDR=127.0.0.1:8080
Environment=OATD_POLICY_PATH=/etc/oadtd/policy.json
WorkingDirectory=/opt/oadtd/current
ExecStart=/opt/oadtd/current/oadtd --addr ${OATD_ADDR} --instance-name %i --public-url ${OATD_PUBLIC_URL} --policy ${OATD_POLICY_PATH} --postgres-dsn ${OATD_POSTGRES_DSN} --api-token ${OATD_API_TOKEN} --alert-webhook-url ${OATD_ALERT_WEBHOOK_URL} --alert-webhook-token ${OATD_ALERT_WEBHOOK_TOKEN} --ticket-webhook-url ${OATD_TICKET_WEBHOOK_URL} --ticket-webhook-token ${OATD_TICKET_WEBHOOK_TOKEN} --response-webhook-url ${OATD_RESPONSE_WEBHOOK_URL} --response-webhook-token ${OATD_RESPONSE_WEBHOOK_TOKEN} --github-api-base ${OATD_GITHUB_API_BASE} --github-owner ${OATD_GITHUB_OWNER} --github-repo ${OATD_GITHUB_REPO} --github-token ${OATD_GITHUB_TOKEN} --github-workflow-file ${OATD_GITHUB_WORKFLOW_FILE} --github-workflow-ref ${OATD_GITHUB_WORKFLOW_REF} --mcp-upstream-url ${OATD_MCP_UPSTREAM_URL} --mcp-upstream-token ${OATD_MCP_UPSTREAM_TOKEN} --oidc-issuer-url ${OATD_OIDC_ISSUER_URL} --oidc-client-id ${OATD_OIDC_CLIENT_ID} --oidc-client-secret ${OATD_OIDC_CLIENT_SECRET} --oidc-redirect-url ${OATD_OIDC_REDIRECT_URL} --oidc-scopes ${OATD_OIDC_SCOPES} --oidc-tenant-claim ${OATD_OIDC_TENANT_CLAIM} --oidc-role-claim ${OATD_OIDC_ROLE_CLAIM} --oidc-email-claim ${OATD_OIDC_EMAIL_CLAIM} --saml-root-url ${OATD_SAML_ROOT_URL} --saml-idp-metadata-url ${OATD_SAML_IDP_METADATA_URL} --saml-key-path ${OATD_SAML_KEY_PATH} --saml-cert-path ${OATD_SAML_CERT_PATH} --saml-tenant-attribute ${OATD_SAML_TENANT_ATTRIBUTE} --saml-role-attribute ${OATD_SAML_ROLE_ATTRIBUTE} --saml-email-attribute ${OATD_SAML_EMAIL_ATTRIBUTE} --trusted-proxies ${OATD_TRUSTED_PROXIES} --retention-window ${OATD_RETENTION_WINDOW} --gateway-max-in-flight ${OATD_GATEWAY_MAX_IN_FLIGHT}
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
