# packaging/systemd/oadtd.service
[Unit]
Description=Open Agentic Threat Defense
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oadtd
Group=oadtd
WorkingDirectory=/opt/oadtd/current
# Bind loopback only; front with a reverse proxy or SSH-tunnel for UI access.
ExecStart=/opt/oadtd/current/oadtd --addr 127.0.0.1:8080 --policy /etc/oadtd/policy.json
EnvironmentFile=-/etc/oadtd/oadtd.env
Restart=on-failure
RestartSec=5s

# Sandbox hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/oadtd
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RemoveIPC=true
UMask=0077
TasksMax=256
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
