# packaging/systemd/oadtd-validate-alert.service
[Unit]
Description=OADTD detection validation failure alert
# Triggered via OnFailure= from oadtd-validate.service. Sends one webhook alert
# (rich, with the failed techniques, when a result file is available; otherwise a
# generic "did not complete" notice). Configure the destination in
# /etc/oadtd/validate.env:
#   OATD_VALIDATE_ALERT_URL=https://hooks.example/detection-regression
#   OATD_VALIDATE_ALERT_TOKEN=optional-bearer-token

[Service]
Type=oneshot
User=oadtd
Group=oadtd
EnvironmentFile=-/etc/oadtd/validate.env
ExecStart=/usr/local/sbin/oadtd-validate-alert

# Sandbox hardening (needs outbound network for the webhook + read of the result)
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RemoveIPC=true
UMask=0077
